# dot_config/systemd/user/home-cleaner.service
[Unit]
Description=Home cleaner
[Service]
ExecStart=%h/.bin/home-cleaner

